# sparse/usr/lib/systemd/system/wlan-up.service
Description=Wlan Up script
After=bluetooth.service

[Service]
ExecStart=/usr/bin/droid/wlan-up.sh

[Install]
WantedBy=network.target

